# containers/fixture/shell-basic/profile.d/10-greeting.sh
# Drop-in profile.d snippet, copied to ~/.config/profile.d/10-greeting.sh
# via the auto-on-directory link entry in shell-basic/Dotfile.yaml. Used by
# verify.sh to confirm the directory cascade produces individual file leaves
# under the target.
export DOTS_TEST_GREETING="hello from shell-basic"
